# skills/openclaw-feishu-multi-agent-deploy/templates/systemd/v4-3-watchdog.service
[Unit]
Description=OpenClaw V4.3.1 single-group watchdog
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env python3 %h/.openclaw/workspace-supervisor_agent/.openclaw/v4_3_job_registry.py --db %h/.openclaw/workspace-supervisor_agent/.openclaw/team_jobs.db watchdog-tick --group-peer-id __TEAM_GROUP_PEER_ID__ --stale-seconds __STALE_SECONDS__
