# skills/openclaw-feishu-multi-agent-deploy/templates/systemd/v4-3-watchdog.timer
[Unit]
Description=Run OpenClaw V4.3.1 single-group watchdog every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Persistent=true
Unit=v4-3-watchdog.service

[Install]
WantedBy=timers.target
