# skills/openclaw-feishu-multi-agent-deploy/templates/systemd/v51-team-watchdog.service
[Unit]
Description=OpenClaw V5.1 team watchdog (__TEAM_KEY__)
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env python3 __RECONCILE_SCRIPT__ --manifest __MANIFEST_PATH__ --team-key __TEAM_KEY__ --openclaw-home __OPENCLAW_HOME__ resume-job --stale-seconds __STALE_SECONDS__
WorkingDirectory=__TEAM_WORKDIR__
Environment=OPENCLAW_SUPERVISOR_AGENT_ID=__SUPERVISOR_AGENT_ID__
