# skills/openclaw-feishu-multi-agent-deploy/templates/systemd/v51-team-watchdog.timer
[Unit]
Description=Run OpenClaw V5.1 team watchdog (__TEAM_KEY__) every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Persistent=true
Unit=v51-team-__TEAM_KEY__.service

[Install]
WantedBy=timers.target
