# skills/openclaw-feishu-multi-agent-deploy/templates/launchd/v5-team-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>bot.molt.v5-team-__TEAM_KEY__</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>__REGISTRY_SCRIPT__</string>
    <string>--db</string>
    <string>__DB_PATH__</string>
    <string>watchdog-tick</string>
    <string>--group-peer-id</string>
    <string>__TEAM_GROUP_PEER_ID__</string>
    <string>--stale-seconds</string>
    <string>__STALE_SECONDS__</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__TEAM_WORKDIR__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/v5-team-__TEAM_KEY__.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/v5-team-__TEAM_KEY__.stderr.log</string>
</dict>
</plist>
